# crons/crontab

* * * * * /home/www-data/repos/smr_web/deployconf/service_reloader